# planroom-genius.service
# Planroom Genius Systemd Service
# Optimized for Raspberry Pi 5
#
# Installation:
#   sudo cp planroom-genius.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable planroom-genius
#   sudo systemctl start planroom-genius
#
# Commands:
#   sudo systemctl status planroom-genius  # Check status
#   sudo systemctl restart planroom-genius # Restart
#   sudo journalctl -u planroom-genius -f  # View logs

[Unit]
Description=Planroom Genius - Construction Lead Intelligence
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as the pi user (change if using different username)
User=pi
Group=pi

# Working directory (change path if installed elsewhere)
WorkingDirectory=/home/pi/planroom-genius

# Start command — uses start_pi.sh which kills old instances before launching
ExecStart=/bin/bash /home/pi/planroom-genius/pi5_build/start_pi.sh

# Restart policy
Restart=always
RestartSec=10

# Environment variables for headless operation
Environment="HEADLESS=true"
Environment="HOME=/home/pi"
Environment="XDG_CONFIG_HOME=/home/pi/.config"
Environment="PLAYWRIGHT_BROWSERS_PATH=/home/pi/.cache/ms-playwright"

# Resource limits (optimized for Pi 5)
# Pi 5 4GB: MemoryMax=2G
# Pi 5 8GB: MemoryMax=4G
MemoryMax=2G
CPUQuota=80%

# Security hardening (relaxed enough for Playwright browser profiles)
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=planroom-genius

# Timeout settings
TimeoutStartSec=90
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
